# grace_addon/files/general/etc/services.d/php-fpm/finish
#!/usr/bin/with-contenv bash

echo "PHP-FPM service stopped"