# grace_addon/files/php83/etc/services.d/php-fpm83/run
#!/usr/bin/with-contenv bash
set -e

echo "Starting PHP-FPM..."

# Ensure PHP-FPM directories exist
mkdir -p /var/log/php83
mkdir -p /run/php-fpm83
mkdir -p /data

# Set proper permissions
chown -R php:php /var/log/php83
chown -R php:php /run/php-fpm83
chown -R php:nginx /data
chown -R php:nginx /var/www

# Ensure database file exists and has correct permissions
if [ -f "/var/www/public/grace.db" ]; then
    chown php:nginx /var/www/public/grace.db
    chmod 664 /var/www/public/grace.db
fi

# Test PHP-FPM configuration
php-fpm83 -t

# Start PHP-FPM
exec php-fpm83 -F -R
